# etc/udev/rules.d/99-disable-capslock.rules
SUBSYSTEM=="input", ACTION=="add", ENV{ID_CLASS}=="kbd", RUN+="/usr/local/sbin/capslock"
